# jupiter-tweaks-desktop/etc/xdg/autostart/setup-kwallet.desktop
[Desktop Entry]
Name=Setup KWallet
Type=Application
Exec=/usr/bin/steamos-setup-kwallet
OnlyShowIn=KDE
