# etc/xdg/autostart/setup-kwallet.desktop
[Desktop Entry]
Name=Setup KWallet
Type=Application
Exec=/usr/bin/holo-setup-kwallet
OnlyShowIn=KDE
